# .dev/docker-env/config/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
autostart=true

[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
